# docker/supervisord.conf
[supervisord]
nodaemon=true
silent=true
logfile=/data/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=5
pidfile=/home/birdnet/run/supervisord.pid
childlogdir=/data/logs

[unix_http_server]
file=/home/birdnet/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///home/birdnet/run/supervisor.sock

[program:php-fpm]
command=/usr/sbin/php-fpm8.5 -F
autostart=true
autorestart=true
priority=10
stdout_logfile=/data/logs/php-fpm.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
redirect_stderr=true

[program:caddy]
command=/usr/bin/caddy run --config /etc/caddy/Caddyfile --adapter caddyfile
autostart=true
autorestart=true
priority=20
stdout_logfile=/data/logs/caddy.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
redirect_stderr=true

[program:recording]
command=/home/birdnet/BirdNET-Pi/birdnet/birdnet_recording.sh
directory=/home/birdnet/BirdNET-Pi/birdnet
autostart=%(ENV_AV_AUDIO_AUTOSTART)s
autorestart=true
startsecs=5
startretries=1000
priority=30
stdout_logfile=/data/logs/recording.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
redirect_stderr=true

[program:analysis]
command=/home/birdnet/BirdNET-Pi/.venv/bin/python3 birdnet_analysis.py
directory=/home/birdnet/BirdNET-Pi/birdnet
autostart=%(ENV_AV_AUDIO_AUTOSTART)s
autorestart=true
priority=30
stdout_logfile=/data/logs/analysis.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
redirect_stderr=true

[program:charts]
command=/home/birdnet/BirdNET-Pi/.venv/bin/python3 daily_plot.py --daemon --sleep 2
directory=/home/birdnet/BirdNET-Pi/birdnet
autostart=true
autorestart=true
priority=40
stdout_logfile=/data/logs/charts.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
redirect_stderr=true

[program:stats]
command=/home/birdnet/BirdNET-Pi/.venv/bin/streamlit run plotly_streamlit.py --browser.gatherUsageStats false --server.address localhost --server.baseUrlPath /stats
directory=/home/birdnet/BirdNET-Pi/birdnet
autostart=true
autorestart=true
priority=40
stdout_logfile=/data/logs/stats.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
redirect_stderr=true

[program:icecast]
command=/usr/bin/icecast2 -c /etc/icecast2/icecast.xml
autostart=true
autorestart=true
priority=50
stdout_logfile=/data/logs/icecast.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
redirect_stderr=true

[program:livestream]
command=/home/birdnet/BirdNET-Pi/birdnet/livestream.sh
directory=/home/birdnet/BirdNET-Pi/birdnet
autostart=%(ENV_AV_AUDIO_AUTOSTART)s
autorestart=true
priority=60
stdout_logfile=/data/logs/livestream.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
redirect_stderr=true
