# avian/forwarding/avian-mqtt.service
[Unit]
Description=AvianVisitors MQTT bridge
After=network-online.target

[Service]
Type=simple
# Edit User= if BirdNET-Pi was installed under a different account.
# The ExecStart path uses an absolute path because systemd's %h
# specifier always expands to /root for system-mode units, ignoring
# User=. WorkingDirectory=~ does respect User=, so we use it.
User=birdnet
WorkingDirectory=~
ExecStart=/usr/bin/python3 avian-mqtt.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
